# packaging/debian/percona-backup-mongodb-agent.files
#!/bin/sh
SOURCE_DIR=$1
BUILD_ROOT=$2
cd $SOURCE_DIR

install -m 0755 -d $BUILD_ROOT/usr/bin
install -m 0755 -d $BUILD_ROOT/etc
install -m 0755 -d $BUILD_ROOT/lib/systemd/system
install -m 0755 bin/pbm-agent $BUILD_ROOT/usr/bin
install -m 0640 etc/pbm-agent.conf $BUILD_ROOT/etc
install -m 0640 etc/pbm-agent-storage.conf $BUILD_ROOT/etc
install -m 0755 lib/systemd/system/pbm-agent.service $BUILD_ROOT/lib/systemd/system